# packaging/mused-server.socket
[Socket]
SocketUser=system
SocketGroup=system
ListenStream=/tmp/.mmsvc_core_socket
Service=mused-server.service

[Install]
WantedBy=sockets.target
