# packaging/muse-server.socket
[Socket]
ListenStream=/tmp/.mmsvc_core_socket
Service=muse-server.service

[Install]
WantedBy=sockets.target
